# http/minikube
#!/sbin/openrc-run

MINIKUBE_BINARY="${MINIKUBE_BINARY:-/usr/bin/minikube}"

name="$RC_SVCNAME"
description="Runs minikube on startup"

depend() {
    need docker
}
start() {
    ebegin "Starting $RC_SVCNAME"
    su alpine -c "$MINIKUBE_BINARY start"
    eend $?
}

stop() {
    ebegin "Stopping $RC_SVCNAME"
    su alpine -c "$MINIKUBE_BINARY stop"
    eend $?
}

status() {
	local _retval
    su alpine -c "$MINIKUBE_BINARY stop" > /dev/null
    _retval=$?
    if [ ${_retval} = '0' ]; then
        einfo 'status: started'
        return 0
    else
        einfo 'status: stopped' 
        return 3
    fi
}
